# laptop/notify/notifier.service
[Service]
Type=simple
ExecStart=/usr/bin/magnotifier.py
