# managed/highstate.service
[Unit]
Description=Run salt on this host

[Service]
Type=simple
ExecStart=/usr/bin/salt-call state.highstate

