# containers/sync/containersync.service
[Unit]
Description=Sync containers among cloud hosts

[Service]
ExecStart=/usr/bin/containersync

[Install]
WantedBy=multi-user.target
