# challenges/challenges.service
[Unit]
Description=Run Magfest Challenges

[Service]
ExecStart=/usr/bin/startchallenges

[Install]
WantedBy=challenges.target
